# kafka/supervisor/zookeeper.conf
[program:zkeeper]
command=/usr/bin/start-zookeeper.sh
autostart=true
stopasgroup=true
